# agents/loki/loki-dashboard.service
[Unit]
Description=Loki Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/nexus-core
Environment="PYTHONPATH=/root/nexus-core"
ExecStart=/root/nexus-core/api/venv/bin/python /root/nexus-core/agents/loki/loki_dashboard.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loki-dashboard

[Install]
WantedBy=multi-user.target
